# systemd/home-hud.service
[Unit]
Description=Home HUD e-ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hud
Group=hud
SupplementaryGroups=audio
WorkingDirectory=/opt/homehud/src
ExecStart=/opt/homehud/venv/bin/python main.py
Restart=on-failure
RestartSec=30
EnvironmentFile=/opt/homehud/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=home-hud

[Install]
WantedBy=multi-user.target