# systemd/home-hud.service
[Unit]
Description=Home HUD e-ink Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hud
Group=hud
WorkingDirectory=/opt/home-hud/src
ExecStart=/opt/home-hud/venv/bin/python main.py
Restart=on-failure
RestartSec=30

# Environment
EnvironmentFile=/opt/home-hud/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=home-hud

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/home-hud/output
PrivateTmp=true

[Install]
WantedBy=multi-user.target
